# chef-solo/site-cookbooks/janek/templates/ruby/client_in_loop.service.erb
[Unit]
Description=Client in loop to stress Shaw's Sinatra app
After=systemd-remount-fs.service
ConditionPathExists=/home/vagrant/sinatra/client_in_loop.rb
## WARNING: This file is being managed by Chef.
## Do not edit this file on a server managed by Chef. All your modifications
## will be lost after a full Chef client run.
## Edit this file directly in the git repository.


[Service]
Type=simple
GuessMainPID=True
User=vagrant
Group=vagrant
WorkingDirectory=<%= @sinatra_home_dir %>
ExecStart=/usr/bin/ruby client_in_loop.rb
TimeoutSec=5

[Install]
WantedBy=multi-user.target
