# chef-solo/site-cookbooks/janek/templates/ruby/sinatra.service.erb
[Unit]
Description=Sinatra for Shaw API
After=systemd-remount-fs.service
ConditionPathExists=/home/vagrant/sinatra/config.ru
ConditionPathExists=/home/vagrant/sinatra/shaw.rb
## WARNING: This file is being managed by Chef.
## Do not edit this file on a server managed by Chef. All your modifications
## will be lost after a full Chef client run.
## Edit this file directly in the git repository.

[Service]
Type=simple
GuessMainPID=True
User=vagrant
Group=vagrant
WorkingDirectory=<%= @sinatra_home_dir %>
ExecStart=/usr/bin/rackup
TimeoutSec=5
Restart=always

[Install]
WantedBy=multi-user.target
