# chef-solo/site-cookbooks/janek/templates/prometheus/node_exporter.service.erb
[Unit]
Description=Node Exporter for Prometheus
Wants=network-online.target
After=network-online.target

[Service]
User=<%= @user %>
Group=<%= @user %>
Type=simple
ExecStart=<%= @bin_dir %>/node_exporter
[Install]
WantedBy=multi-user.target
