# chef-solo/site-cookbooks/janek/templates/prometheus/node_exporter.service.erb
[Unit]
Description=Node Exporter for Prometheus
Wants=network-online.target
After=network-online.target
## WARNING: This file is being managed by Chef.
## Do not edit this file on a server managed by Chef. All your modifications
## will be lost after a full Chef client run.
## Edit this file directly in the git repository.

[Service]
User=<%= @user %>
Group=<%= @user %>
Type=simple
ExecStart=<%= @bin_dir %>/node_exporter
[Install]
WantedBy=multi-user.target
