# chef-solo/site-cookbooks/janek/templates/prometheus/prometheus.service.erb
[Unit]
Description=Prometheus
Wants=network-online.target
After=network-online.target
## WARNING: This file is being managed by Chef.
## Do not edit this file on a server managed by Chef. All your modifications
## will be lost after a full Chef client run.
## Edit this file directly in the git repository.


[Service]
User=<%= @user %>
Group=<%= @user %>
Type=simple
GuessMainPID=yes
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=<%= @bin_dir %>/prometheus \
    --config.file <%= @config_dir %>/prometheus.yml \
    --storage.tsdb.path <%= @var_dir %>/ \
    --web.console.templates=<%= @config_dir %>/consoles \
    --web.console.libraries=<%= @config_dir %>/console_libraries

[Install]
WantedBy=multi-user.target
