# chef-solo/site-cookbooks/janek/templates/mtail/mtail.service.erb
[Unit]
Description=mtail extract whitebox monitoring data from application logs for collection into a timeseries database
After=syslog.target network.target
## WARNING: This file is being managed by Chef.
## Do not edit this file on a server managed by Chef. All your information will be lost
## after a full Chef client run.
## Edit this file directly in the git repository.

[Service]
Type=simple
GuessMainPID=True
User=root
Group=root
WorkingDirectory=/vagrant/
ExecStart=<%= @bin_file %> -port 3093 -log_dir <%= @log_dir %> -logs <%= @log_dir %>access.log -address 0.0.0.0 -metric_push_interval_seconds 5 -progs <%= @mtail_config_file %>
TimeoutSec=5
Restart=always

[Install]
WantedBy=multi-user.target
